# ic880a-backplane/rootfs-overlay/etc/profile.d/motd.sh
echo ''
echo '=*=*= Welcome to the ic880a backplane test image! =*=*='
echo ''
echo 'Simply use the `t` command to start the tests.'
echo ''
