# ic880a-backplane/rootfs-overlay/etc/profile.d/ps1.sh
BRed='\e[1;31m'         # Red
Color_Off='\e[0m'       # Text Reset

export PS1="$BRed[\w]#$Color_Off "
